# app/supervisord.conf
[supervisord]
nodaemon=true

[supervisorctl]

[program:bind9]
command=/usr/sbin/named -f -u bind
user=root
autostart=false